# configs/supervisor-brain-api.conf
[program:brain-api]
directory=/opt/ai-genesis/brain-api
command=/usr/local/bin/node /opt/ai-genesis/brain-api/server.js
user=ubuntu
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/tmp/brain-api.log
stderr_logfile=/tmp/brain-api.err.log
environment=DATA_DIR="/data",BRAIN_CONFIG_PATH="/opt/ai-genesis/brain-config.json",KNOWLEDGE_DIR="/data/knowledge",THINK_SCRIPT="/opt/scripts/think.sh",PORT="8080"
